# cookbooks/endeca/templates/default/workbench_service.erb
#!/bin/sh
#
# Oracle Endeca Tools Service
#
# chkconfig: 35 95 20
# description: Oracle Endeca Tools Service (Workbench)
#

ENDECA_TOOLS_ROOT=<%= node[:endeca][:installation_dir] %>/endeca/ToolsAndFrameworks/3.1.2
ENDECA_USER=<%= node[:endeca][:user] %>

usage() {
  echo "Usage: ${0} (start|stop)"
  echo ""
  echo "  start   - starts Workbench as a background process"
  echo "  stop    - stops Workbench"
  echo ""
}

case "${1}" in
  start)
    /bin/su - "${ENDECA_USER}" -c "${ENDECA_TOOLS_ROOT}/server/bin/startup.sh"
    ;;
  stop)
    /bin/su - "${ENDECA_USER}" -c "${ENDECA_TOOLS_ROOT}/server/bin/shutdown.sh"
    ;;
  *)
    usage
    exit 2
esac

exit $?